# daemon/smb4.mount.in
[Mount]
Type=smb4-share
Exec=@libexecdir@/gvfsd-smb4
AutoMount=false
Scheme=smb4
